# docker/crontab
* * * * *  php /code/polling.php >> /dev/stdout 2>&1
* * * * *  echo "Hello world" >> /dev/stdout 2>&1
